# pkg/systemd/user/sichter-mcp.service
[Unit]
Description=Sichter MCP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/sichter/mcp/server.js
Restart=always

[Install]
WantedBy=default.target
