# pkg/systemd/sichter-sweep.service
[Unit]
Description=Sichter PR Sweep (oneshot)

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc '%h/sichter/cli/sweep --changed'
