# pkg/systemd/sichter-sweep.timer
[Unit]
Description=Run Sichter PR Sweep every 12h

[Timer]
OnBootSec=5m
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=default.target
